# root/etc/services.d/gunicorn/run
#!/usr/bin/with-contenv bash

exec /usr/bin/gunicorn --config /etc/gunicorn/conf.d/gunicorn.ini app:app
